# .docker/conf.d/metrics.conf
[program:metrics]
command=php -S 0.0.0.0:80 -t pub
numprocs=1
directory=/project
autostart=true
autorestart=true
