# .docker/conf.d/monitor.conf
[program:monitor]
command=php console.php demo:monitor
numprocs=1
directory=/project
autostart=true
autorestart=true
