# .docker/conf.d/worker.conf
[program:worker]
process_name=%(program_name)s_%(process_num)02d
command=php console.php demo:worker
numprocs=20
directory=/project
autostart=false
autorestart=true
stopwaitsecs=20
