# .docker/conf.d/zombiekiller.conf
[program:zombie-killer]
command=php console.php demo:zombie-killer
numprocs=1
directory=/project
autostart=true
autorestart=true
